# taiga/files/supervisor.conf
{%- from "taiga/map.jinja" import server with context %}
{%- set service = salt['pillar.get']('supervisor:server:service:'+service_name) %}

[program:{{ service_name }}]
environment=PYTHONOPTIMIZE=1,PATH="{{ server.virtualenv }}/bin"
directory={{ server.dir }}/taiga-back

command={{ server.virtualenv }}/bin/gunicorn -w 3 -t 60 --pythonpath=. -b 0.0.0.0:8001 taiga.wsgi

stdout_logfile=/var/log/taiga/{{ service.name }}_stdout.log
stderr_logfile=/var/log/taiga/{{ service.name }}_stderr.log

user=taiga
autostart=true
autorestart=true

[program:{{ service_name }}-celery]
environment=PYTHONOPTIMIZE=1,PATH="{{ server.virtualenv }}/bin"
directory={{ server.dir }}/taiga-back

command={{ server.virtualenv }}/bin/celery -A taiga worker -c 4

stdout_logfile=/var/log/taiga/{{ service.name }}-celery_stdout.log
stderr_logfile=/var/log/taiga/{{ service.name }}-celery_stderr.log

user=taiga
autostart=true
autorestart=true
